# oooo.conf
[program:oooo]
command=python3 /home/curator2/oooo/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/oooo/test.err.log
stdout_logfile=/home/curator2/oooo/test.out.log